# releases/sources/qnapware/opt/lib/opkg/info/qnapware.list
/Apps/opt/etc/init.d/rc.func
/Apps/opt/etc/init.d/rc.unslung
/Apps/opt/etc/profile
